# services/mc-idle-check.timer
[Unit]
Description=Runs Minecraft Idle Check every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target

